# worker/com.bookedjob.podcast.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.bookedjob.podcast</string>
  <key>ProgramArguments</key><array>
    <string>/bin/bash</string><string>-lc</string>
    <string>cd /Users/aaronphillips/GIT/booked-job && /usr/bin/python3 scripts/podcast_build.py && git add -A && git commit -q -m "podcast: new episode" && git push -q</string></array>
  <key>WorkingDirectory</key><string>/Users/aaronphillips/GIT/booked-job</string>
  <key>StartCalendarInterval</key><dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
  <key>StandardOutPath</key><string>/Users/aaronphillips/GIT/booked-job/content/podcast.out.log</string>
  <key>StandardErrorPath</key><string>/Users/aaronphillips/GIT/booked-job/content/podcast.err.log</string>
  <key>RunAtLoad</key><false/>
</dict></plist>
