# worker/com.bookedjob.publisher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bookedjob.publisher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/aaronphillips/GIT/booked-job/scripts/publisher.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/aaronphillips/GIT/booked-job</string>

  <!-- Fire inside each posting window; publisher's per-day chosen slot decides
       which firing actually posts (one post/day max during warm-up). -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>10</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>40</integer></dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/aaronphillips/GIT/booked-job/content/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/aaronphillips/GIT/booked-job/content/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
